# archiso/airootfs/etc/systemd/system/applesmc_load.service
# Load applesmc module early for fan control
[Unit]
Description=Load applesmc sensor module

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/modprobe applesmc

[Install]
WantedBy=multi-user.target